# starling_setup/voxl-aeroprint-start.service
Description=Aeroprint Startup Service

[Service]
ExecStart=/etc/aeroprint/startup.bash

[Install]
WantedBy=multi-user.target
